# audio-server-binary.service
[Unit]
Description=HTTP Audio Server (Binary)
After=network.target sound.target

[Service]
Type=simple
User=alma
Group=alma
Environment="DISPLAY=:0"
Environment="PULSE_SERVER=unix:/run/user/1000/pulse/native"
ExecStart=/usr/local/bin/audio-server --host 0.0.0.0 --port 8800
Restart=on-failure
RestartSec=5s

# Audio access
SupplementaryGroups=audio

[Install]
WantedBy=multi-user.target
